# capslock-auto-switch/etc/systemd/user/capslock-auto-switch.service
[Unit]
Description=Auto switch capslock

[Service]
Type=simple
ExecStart=/usr/bin/capslock-auto-switch

[Install]
WantedBy=multi-user.target
